# ansible/build/roles/install-jenkins/files/alpha-login.service
[Unit]
Description=Alpha Login Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/login-job.sh
StandardOutput=journal
StandardError=journal
